# debian/postinst
#!/bin/sh

update-rc.d wb-homa-ism-radio defaults
invoke-rc.d wb-homa-ism-radio start
